# althea/althea-rust-binaries/files/rita.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2006-2016 OpenWrt.org

START=90
PROG=/usr/sbin/rita

config_file=/etc/rita.toml
log_file=/var/log/rita.log
RUST_LOG=INFO

start() {
	[ -f $config_file ] || return 1
	sleep 30
	(service_start $PROG --config=/etc/rita.toml --platform=linux &> $log_file)&
}

stop() {
	service_stop $PROG
}

reload() {
	$PROG2 reload
}
